# arnold_pipelines/megaplan/cloud/systemd/megaplan-watchdog-ensure.timer
[Unit]
Description=Ensure Megaplan cloud watchdog is alive every minute

[Timer]
OnBootSec=1min
OnUnitActiveSec=1min
AccuracySec=5s
Persistent=true
Unit=megaplan-watchdog-ensure.service

[Install]
WantedBy=timers.target
