# agentbox/systemd/agentbox-discord-resident.service
[Unit]
Description=AgentBox Discord resident
After=network.target

[Service]
Type=simple
WorkingDirectory={{WORKING_DIRECTORY}}
Environment=AGENTBOX_CONFIG={{WORKSPACE_ROOT}}/agentbox.yaml
Environment=PYTHONPATH={{WORKING_DIRECTORY}}
# Set the real token with a drop-in, e.g.
#   sudo systemctl edit agentbox-discord-resident
# and add:
#   [Service]
#   Environment=DISCORD_BOT_TOKEN=<token>
Environment=DISCORD_BOT_TOKEN=
ExecStart={{PYTHON}} -m arnold_pipelines.megaplan resident discord
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
