# agentbox/systemd/arnold-guardian.service
[Unit]
Description=Arnold Guardian worker for AgentBox
After=network.target

[Service]
Type=simple
WorkingDirectory={{WORKING_DIRECTORY}}
Environment=AGENTBOX_CONFIG={{WORKSPACE_ROOT}}/agentbox.yaml
Environment=PYTHONPATH={{WORKING_DIRECTORY}}
ExecStart={{PYTHON}} -m agentbox guardian run
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
